# disksrc/liabarmdist_base/etc/inittab
#
# inittab       This file describes how the INIT process should set up
#               the system in a certain run-level.
#
# Author:       Mikael Dich, <midi@liab.dk>
#
# Default runlevel. The runlevels used by RHS are:
#   0 - halt (Do NOT set initdefault to this)
#   1 - Single user mode
#   2 - Multiuser, without NFS (The same as 3, if you do not have networking)
#   3 - Full multiuser mode
#   4 - unused
#   5 - X11
#   6 - reboot (Do NOT set initdefault to this)
# 
id:3:initdefault:

# System initialization 
si::sysinit:/etc/rc.d/rc.sysinit

# Script to run when going single user:
su:1:wait:/etc/rc.d/rc.S
# Script to run when going multi user:
mu:23456:wait:/etc/rc.d/rc.M


s1:2345:respawn:/sbin/getty -L ttyS0 115200 vt100

